# tasks/crontab
* * * * * /usr/bin/python3 /app/task.py >> /var/log/cron.log 2>&1
